# dev-db/mysql/files/mysql.init
#!/sbin/openrc-run
# Copyright 2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

MYSQLD="/usr/@LIBDIR@/mysql-@SLOT@/sbin/mysqld"

description="MySQL @SLOT@ -- Fast, multi-threaded, multi-user SQL database server"

depend() {
    use net
    provide mysql
}
